# crontab
0,50 * * * * /home/omen/CONKY/eloadrssfeeds.sh
